# robot/runrobot.timer
[Unit]
Description=Load network drive with delay

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target

